# systemd/asterisk.service.d/cpu-scheduling-policy.conf
[Service]
CPUSchedulingPolicy=rr
CPUSchedulingPriority=95
