# systemd/posgresql.service.d/cpu-scheduling-policy.conf
[Service]
Nice=15
CPUSchedulingPolicy=idle
